# reference/rc.local
#!/usr/bin/env sh

## prevent immediate wakeup due to USB, ethernet on X250

echo EHC1 > /proc/acpi/wakeup
echo XHCI > /proc/acpi/wakeup
echo IGBE > /proc/acpi/wakeup
